# Code/Chapter09/ubuntu/etc/crontab
0 */2 * * *    root    /usr/sbin/ntpdate pool.ntp.org > /dev/null 2>&1
